# site-env/files/site.conf
[program:site]
command=/vagrant/files/uwsgi --ini /vagrant/files/uwsgi.ini
directory=/vagrant/dmoj/site
stopsignal=QUIT
stdout_logfile=/tmp/site.stdout.log
stderr_logfile=/tmp/site.stderr.log
