# site-env/files/site.conf
[program:site]
command=/vagrant/files/uwsgi --ini /vagrant/files/uwsgi.ini
directory=/vagrant/site
stopsignal=QUIT
stdout_logfile=/vagrant/logs/site.log
stderr_logfile=/vagrant/logs/site.log
