# judge-env/systemd_files/dmoj-sync-problems.timer
[Unit]
Description=DMOJ sync problems timer

[Timer]
OnActiveSec=0
OnUnitInactiveSec=15min
AccuracySec=3min

[Install]
WantedBy=multi-user.target
