# judge-env/systemd_files/dmoj-judge@.service
[Unit]
Description=DMOJ judge
After=network.target

[Service]
Type=simple
User=nobody
ExecStart=/usr/local/bin/dmoj -c /vagrant/judge.yml 52.28.196.151 %i
Restart=always

[Install]
WantedBy=multi-user.target
