# judge-env/systemd_files/dmoj-sync-problems.service
[Unit]
Description=DMOJ sync problems
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/problems
ExecStart=/usr/bin/git pull
